# com.jpothen.powermonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.jpothen.powermonitor</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/jpothen/Developer/Power_Info/.venv/bin/power_monitor</string>
		<string>/Users/jpothen/Developer/Power_Info/power_monitor.py</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/jpothen/Developer/Power_Info</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<false/>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>StandardOutPath</key>
	<string>/Users/jpothen/Developer/Power_Info/logs/power_monitor.out.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/jpothen/Developer/Power_Info/logs/power_monitor.err.log</string>
</dict>
</plist>
